# scripts/system/kiro-gateway.service
[Unit]
Description=Kiro Gateway (proxy for Kiro API)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=${INSTALL_DIR}/state
ExecStart=${INSTALL_DIR}/venv/bin/python main.py
Environment=KIRO_GATEWAY_HOME=${INSTALL_DIR}
Environment=ACCOUNTS_CONFIG_FILE=${INSTALL_DIR}/state/credentials.json
Environment=ACCOUNTS_STATE_FILE=${INSTALL_DIR}/state/state.json
StandardOutput=append:${INSTALL_DIR}/logs/kiro-gateway.log
StandardError=append:${INSTALL_DIR}/logs/kiro-gateway.log

[Install]
WantedBy=default.target
