# scripts/system/kiro-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.jwadow.kiro-gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>${INSTALL_DIR}/venv/bin/python</string>
    <string>main.py</string>
  </array>
  <key>WorkingDirectory</key> <string>${INSTALL_DIR}/app</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>KIRO_GATEWAY_HOME</key>     <string>${INSTALL_DIR}</string>
    <key>ACCOUNTS_CONFIG_FILE</key>  <string>${INSTALL_DIR}/state/credentials.json</string>
    <key>ACCOUNTS_STATE_FILE</key>   <string>${INSTALL_DIR}/state/state.json</string>
  </dict>
  <key>RunAtLoad</key>      <false/>
  <key>KeepAlive</key>      <false/>
  <key>StandardOutPath</key> <string>${INSTALL_DIR}/logs/kiro-gateway.log</string>
  <key>StandardErrorPath</key><string>${INSTALL_DIR}/logs/kiro-gateway.log</string>
  <key>ProcessType</key>    <string>Background</string>
</dict>
</plist>
